# lgh_docker/etc/supervisor/logstash.ini
[program:logstash]
command = /home/logstash/logstash-6.2.3/bin/logstash -f /home/logstash/logstash-6.2.3/config/logstash-filebeat.yml
user = logstash
autorestart = true
stderr_logfile = /tmp/logstash-supervisor.err
stdout_logfile = /tmp/logstash-supervisor.log
