# kib_docker/etc/supervisor/kibana.ini
[program:kibana]
command = /home/kibana/kibana-6.2.3-linux-x86_64/bin/kibana
user = kibana
autorestart = true
stderr_logfile = /tmp/kibana-supervisor.err
stdout_logfile = /tmp/kibana-supervisor.log
