# es_docker/etc/supervisor/elasticsearch.ini
[program:elasticsearch]
command = /home/es/elasticsearch-6.2.3/bin/elasticsearch
user = es
autorestart = true
stderr_logfile = /tmp/elasticsearch-supervisor.err
stdout_logfile = /tmp/elasticsearch-supervisor.log
