# es_docker/etc/supervisor/ntpdate.ini
[program:ntpdate]
command = bash -c 'while true; do echo ntpdate 0.pool.ntp.org;ntpdate 0.pool.ntp.org;sleep 60;done ;'
autorestart = true
stderr_logfile = NONE
stdout_logfile = NONE
